# getmail/.config/systemd/user/getmail.timer
[Unit]
Description=getmail timer

[Timer]
OnCalendar=*:0/30
Persistent=true
Unit=getmail.service

[Install]
WantedBy=timers.target
