# getmail/.config/systemd/user/getmail.service
[Unit]
Description=getmail service

[Service]
Type=oneshot
ExecStart=/usr/bin/getmail
StandardOutput=syslog
StandardError=syslog
