# wp/.config/systemd/user/wp.service
[Unit]
Description=wallpaper service

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c %h/.local/bin/wp-auto
