# notmuch/.config/systemd/user/notmuch.service
[Unit]
Description=notmuch service

[Service]
Type=oneshot
ExecStart=/usr/bin/notmuch new
StandardOutput=syslog
StandardError=syslog
